# etc/supervisord.d/rsyslogd.ini
[program:rsyslogd]
command=/usr/sbin/rsyslogd -n

# Redirect logs to stdout
stdout_events_enabled=true
stderr_events_enabled=true
